# solarbank_load_manager/rootfs/etc/services.d/solarbank-manager/run
#!/usr/bin/with-contenv sh
set -e
cd /app/backend
exec python3 -m solarbank_manager.main
